# scripts/gpu-monitor.service
[Unit]
Description=GPU Monitor — DCGM exporter 기반 다중 GPU VM 통합 모니터링
After=network.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
ExecStart=__VENV_DIR__/bin/gpu-monitor -c __INSTALL_DIR__/config.yaml start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
